# checkbattery.service
[Unit]
Description=bateria

[Service]
Restart=always
RestartSec=60
ExecStart=check-battery

[Install]
WantedBy=multi-user.target

